# volumes/supervisor/conf.d/backup-service.conf
[program:backup-service]
command=gunicorn --workers=3 manage:app -b 0.0.0.0:5000
directory=/usr/src
autorestart=true
startretries=50

[supervisord]
nodaemon=true
